# backupmanager/systemd/backup.timer
[Unit]
Description=Backupmanager timer

[Timer]
OnCalendar=daily
Unit=backup.service

[Install]
WantedBy=multi-user.target
