# backupmanager/systemd/backup.service
[Unit]
Description=Backupmanager service

[Service]
Type=simple
ExecStart=/usr/local/bin/backupmanager run

[Install]
WantedBy=backup.timer